# ops/systemd/computor-analytics-source-tunnel.service
[Unit]
Description=Computor analytics source database tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/computor/analytics-source-tunnel.env
Environment=ANALYTICS_SOURCE_TUNNEL_ENV_FILE=/etc/computor/analytics-source-tunnel.env
ExecStart=/bin/bash -lc 'cd "$COMPUTOR_DEPLOY_DIR" && exec scripts/analytics-prod/source-tunnel.sh foreground'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
